# monitoring.service
[Unit]
Description=LegacyCoinTrader System Service
After=network.target docker.service
Wants=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=brandonburnette
WorkingDirectory=/Users/brandonburnette/Downloads/LegacyCoinTrader1.0-1
ExecStart=/bin/bash /Users/brandonburnette/Downloads/LegacyCoinTrader1.0-1/startup.sh start
ExecStop=/bin/bash /Users/brandonburnette/Downloads/LegacyCoinTrader1.0-1/shutdown.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/Users/brandonburnette/Downloads/LegacyCoinTrader1.0-1
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/Users/brandonburnette/Downloads/LegacyCoinTrader1.0-1
PrivateTmp=yes

# Resource limits
MemoryLimit=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
